# Denis_Version/systemd/system/camera_webrtc.service
# ============================================================================
# Systemd Service für Raspberry Pi Kamera-System (WebRTC-Version)
# Datei: /etc/systemd/system/camera_webrtc.service
# Version: 1.0
# Datum: 30.09.2025
#
# Installation:
# sudo cp camera_webrtc.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable camera_webrtc.service
# sudo systemctl start camera_webrtc.service
#
# Befehle:
# Status prüfen:    sudo systemctl status camera_webrtc.service
# Starten:          sudo systemctl start camera_webrtc.service
# Stoppen:          sudo systemctl stop camera_webrtc.service
# Neu starten:      sudo systemctl restart camera_webrtc.service
# Logs anzeigen:    sudo journalctl -u camera_webrtc.service -f
# ============================================================================

[Unit]
# Beschreibung des Services
Description=Raspberry Pi 5 Kamera-Streaming-System (WebRTC)
Documentation=https://github.com/raspberrypi/camera-system

# Wartet bis Netzwerk verfügbar ist
After=network-online.target
Wants=network-online.target

# Startet nach Multiuser-System
After=multi-user.target

[Service]
# Service-Typ
Type=simple

# Benutzer unter dem der Service läuft
User=pi
Group=pi

# Arbeitsverzeichnis
WorkingDirectory=/home/pi/camera_system

# Python Virtual Environment aktivieren und Hauptprogramm starten
ExecStart=/home/pi/camera_system/venv/bin/python3 /home/pi/camera_system/main_webrtc.py

# Automatischer Neustart bei Absturz
Restart=always
RestartSec=10

# Standard-Output und Fehler-Output an Systemd-Journal
StandardOutput=journal
StandardError=journal

# Syslog-Identifier für bessere Log-Filterung
SyslogIdentifier=camera_webrtc

# Umgebungsvariablen
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/pi/camera_system"

# WebRTC benötigt eventuell mehr Ressourcen
# LimitNOFILE=65536
# MemoryLimit=1.5G
# CPUQuota=90%

# Timeout-Einstellungen
TimeoutStartSec=60
TimeoutStopSec=30

# Kill-Modus
KillMode=mixed
KillSignal=SIGTERM

# Cleanup bei Service-Stop
ExecStopPost=/bin/bash -c 'sleep 2'

[Install]
# Startet automatisch im Multi-User-Modus
WantedBy=multi-user.target
